# deploy/systemd/mosaicfs.service
[Unit]
Description=MosaicFS agent
Documentation=https://github.com/...
After=network-online.target couchdb.service
Wants=network-online.target

[Service]
Type=simple
User=mosaicfs
Group=mosaicfs
ExecStart=/usr/local/bin/mosaicfs --config /etc/mosaicfs/mosaicfs.toml
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=3

StateDirectory=mosaicfs
RuntimeDirectory=mosaicfs
LogsDirectory=mosaicfs
Environment="XDG_RUNTIME_DIR=%t/mosaicfs"

# --- privilege ---
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=

# --- filesystem ---
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ReadWritePaths=/var/lib/mosaicfs /run/mosaicfs

# --- proc/kernel ---
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# --- misc ---
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RemoveIPC=yes
MemoryDenyWriteExecute=yes
UMask=0027

# --- resource caps ---
MemoryMax=512M
CPUQuota=75%
TasksMax=256

# --- network ---
IPAddressDeny=any
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

[Install]
WantedBy=multi-user.target
